# config_templates/gunicorn.service.in
[Unit]
Description=gunicorn daemon
Requires=gunicorn.socket
After=network.target
[Service]
Type=notify
User=admin
Group=www-data
RuntimeDirectory=gunicorn
WorkingDirectory=workdir
Environment="DJANGO_ENV=production"
ExecStart=venvpath/bin/gunicorn \
          --access-logfile - \
          --workers 3 \
          --bind unix:/run/gunicorn.sock \
          server.wsgi:application
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true

[Install]
WantedBy=multi-user.target
